# deploy/templates/ysxs-gunicorn.service.template
[Unit]
Description=YSXS Gunicorn Service
After=network.target

[Service]
User=__RUN_USER__
Group=__RUN_GROUP__
WorkingDirectory=__PROJECT_DIR__
EnvironmentFile=__ENV_FILE__
ExecStart=__VENV_DIR__/bin/gunicorn -w 2 -b 127.0.0.1:__PORT__ YSXS.app:app
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
